# imageroot/systemd/user/api-moduled.service.d/api-moduled-overwrite.conf
[Service]
ExecStartPre=runagent write-user-portal-config
